# home/private_dot_config/systemd/user/update-policy.timer
[Unit]
Description=Run update-policy.service daily

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
